# DEBIAN/postinst
#!/bin/bash

systemctl enable webdis.service
systemctl start webdis.service