# DEBIAN/postinst
#!/bin/bash

update-rc.d webdis enable
service webdis start